# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/local/bin/gunicorn needle_app:app -b :5000 -w 2
directory=/needle-masters

[program:nginx]
command=/usr/sbin/nginx
